# .config/systemd/user/save-notes.service
[Unit]
Description=Commit and push notes repository if changed

[Service]
Type=simple
ExecStart=/home/raoul/bin/commit-and-push.sh /home/raoul/Notes
